# serverfiles/crontab
*/1 *   * * *   root    /root/kouteikanri/venv/bin/python /root/kouteikanri/serverfiles/csvcopy.py
*/2 *   * * *   root    sh /root/kouteikanri/serverfiles/yobi.sh
