# serverfiles/crontab
*/1 *   * * *   root    /root/kouteikanri/venv/bin/python /root/csvcopy.py
*/2 *   * * *   root    sh /root/yobi.sh
